# config/omarchy/systemd/omarchy-theme-schedule.service
[Unit]
Description=Apply scheduled Omarchy theme
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/omarchy-theme-schedule
